# src/teamspeak/build/modules/teamspeak/files/etc/supervisor/conf.d/logs.conf
[program:logs]
command = multitail --mergeall -q 1 '/usr/local/src/teamspeak3-server/logs/*'
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
