# src/teamspeak/build/modules/teamspeak/files/etc/supervisor/conf.d/teamspeak.conf
[program:teamspeak]
command = /usr/local/src/teamspeak3-server/ts3server_minimal_runscript.sh
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
