# src/teamspeak/build/modules/teamspeak/files/etc/supervisor/conf.d/teamspeak.conf
[program:teamspeak]
command = /usr/local/src/teamspeak3-server/ts3server_minimal_runscript.sh
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
