# src/teamspeak/build/modules/teamspeak/files/etc/supervisor/conf.d/teamspeak_stdout.conf
[program:teamspeak_stdout]
command = /usr/bin/tail -F /usr/local/src/teamspeak3-server/logs/*.log
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
